# .local/src/bootc/arch/rootfs/usr/lib/systemd/system/systemd-sysext.service.d/mutable.conf
[Service]
ExecStart=
ExecStart=systemd-sysext refresh --mutable=auto
ExecReload=
ExecReload=systemd-sysext refresh --mutable=auto
